# etc/systemd/system/tree-sinker@.service
[Unit]
Description=Update the %i squashed portage tree
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=tree-sinker %i
